# myTunes/DEBIAN/postinst
#!/bin/bash
PLIST=/System/Library/LaunchDaemons/net.iWatcher.myTunes.plist
/usr/sbin/chown root:wheel $PLIST
/bin/launchctl load $PLIST
