# myTunes/DEBIAN/postinst
#!/bin/bash
/bin/launchctl load /System/Library/LaunchDaemons/net.iWatcher.myTunes.plist
